# crontab
* * * * * DISPLAY=:0.0 feh --randomize --bg-fill /home/mitchyman/Pictures/*
